# etc/rc.d/init.d/perdition.rh.in
#!/bin/sh
######################################################################
# Version:                                              September 1999
# Author:                                   Horms <horms@verge.net.au>
#
# Red Hat init script for perdition
# Adapted from: Tomasz K³oczko, <kloczek@rudy.mif.pg.gda.pl>
#               Michele Marziani <marziani@fe.infn.it>
#
# perdition:
# chkconfig: 2345 81 29
# description: Perdition is a POP3 and IMAP4 proxy designed to \
# enable users to contact a content free pop server that \
# will direct them to a real pop server on a per user basis\
# processname: perdition
#
# perdition
# Mail retrieval proxy server
# Copyright (C) 1999-2003  Horms
# 
# This program is free software; you can redistribute it and/or
# modify it under the terms of the GNU General Public License as
# published by the Free Software Foundation; either version 2 of the
# License, or (at your option) any later version.
# 
# This program is distributed in the hope that it will be useful, but
# WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the GNU
# General Public License for more details.
# 
# You should have received a copy of the GNU General Public License
# along with this program; if not, write to the Free Software
# Foundation, Inc., 59 Temple Place, Suite 330, Boston, MA
# 02111-1307  USA
#
######################################################################

# Source function library.
. /etc/rc.d/init.d/functions

if [ -f /etc/sysconfig/perdition ]; then
  . /etc/sysconfig/perdition
fi

# Please do not edit the values below.
# Rather, please edit /etc/default/perdition
RUN_PERDITION="${RUN_PERDITION:=no}"
FLAGS="${FLAGS:=}"
POP3="${POP3:=yes}"
POP3_FLAGS="${POP3_FLAGS:=}"
POP3S="${POP3S:=yes}"
POP3S_FLAGS="${POP3S_FLAGS:=}"
IMAP4_FLAGS="${IMAP4:=yes}"
IMAP4="${IMAP4_FLAGS:=}"
IMAP4S="${IMAP4S:=yes}"
IMAP4S_FLAGS="${IMAP4S_FLAGS:=}"

if [ "$RUN_PERDITION" != "yes" ]; then
	exit 0
fi

# See how we were called.
case "$1" in
  start)
	if [ "$POP3" = "yes" ]; then
        	echo -n "Starting perdition services (POP3): "
		daemon @prefix@/sbin/perdition.pop3 $FLAGS $POP3_FLAGS
		touch /var/lock/subsys/perdition.pop3
        	echo
	fi

	if [ "$POP3S" = "yes" ]; then
        	echo -n "Starting perdition services (POP3S): "
		daemon @prefix@/sbin/perdition.pop3s $FLAGS $POP3S_FLAGS
		touch /var/lock/subsys/perdition.pop3s
        	echo
	fi

	if [ "$IMAP4" = "yes" ]; then
        	echo -n "Starting perdition services (IMAP4): "
		daemon @prefix@/sbin/perdition.imap4 $FLAGS $IMAP4_FLAGS
		touch /var/lock/subsys/perdition.imap4
		echo
	fi

	if [ "$IMAP4S" = "yes" ]; then
        	echo -n "Starting perdition services (IMAP4S): "
		daemon @prefix@/sbin/perdition.imaps $FLAGS $IMAP4S_FLAGS
		touch /var/lock/subsys/perdition.imaps
		echo
	fi
	;;
  status)
        if [ "$POP3" = "yes" ]; then
        	status perdition.pop3
	fi
        if [ "$POP3S" = "yes" ]; then
        	status perdition.pop3s
	fi
	if [ "$IMAP4" = "yes" ]; then
        	status perdition.imap4
	fi
	if [ "$IMAP4S" = "yes" ]; then
        	status perdition.imaps
	fi
	;;
  restart|reload)
        $0 stop
	$0 start
        ;;
  stop)
	if [ "$POP3" = "yes" ]; then
        	echo -n "Shutting down perdition services (POP3): "
		killproc perdition.pop3
		rm -f /var/lock/subsys/perdition.pop3
        	echo
	fi

	if [ "$POP3S" = "yes" ]; then
        	echo -n "Shutting down perdition services (POP3S): "
		killproc perdition.pop3s
		rm -f /var/lock/subsys/perdition.pop3s
        	echo
	fi

	if [ "$IMAP4" = "yes" ]; then
        	echo -n "Shutting down perdition services (IMAP4): "
		killproc perdition.imap4
		rm -f /var/lock/subsys/perdition.imap4
		echo
	fi

	if [ "$IMAP4S" = "yes" ]; then
        	echo -n "Shutting down perdition services (IMAP4S): "
		killproc perdition.imaps
		rm -f /var/lock/subsys/perdition.imaps
		echo
	fi
	;;
  *)
	echo "Usage: perdition {start|stop|status|restart|reload}"
	exit 1
esac

exit 0
