# debian/perdition.conffiles
/etc/init.d/perdition
/etc/perdition/perdition.conf
/etc/pam.d/perdition
