# debian/perdition-bdb.prerm
#! /bin/sh
# prerm script for perdition

#DEBHELPER#

if [ "$1" = "purge"  -o "$1" = "remove" ]; then
	if [ -f /etc/init.d/perdition ]; then
		invoke-rc.d perdition stop
	fi
	if [ -f /etc/init.d/perdition ]; then
		invoke-rc.d perdition start
	fi
fi

exit 0
