# debian/perdition.postinst
#! /bin/sh
# postinst script for perdition

#DEBHELPER#

set -e

if [ "$1" = "configure" ]; then
	if [ -f /etc/init.d/perdition ]; then
		invoke-rc.d perdition stop
	fi
	update-rc.d perdition defaults >/dev/null
	if [ -f /etc/init.d/perdition ]; then
		invoke-rc.d perdition start
	fi
fi

exit 0
