# debian/perdition.files
etc/pam.d/perdition
etc/perdition/perdition.conf
etc/init.d/perdition
usr/share/doc/perdition/CODING_LOCATIONS
usr/share/doc/perdition/ChangeLog
usr/share/doc/perdition/NEWS
usr/share/doc/perdition/README
usr/share/doc/perdition/README.perditiondb
usr/share/doc/perdition/TODO
usr/lib/libjain.so
usr/lib/libjain.so.0
usr/sbin/perdition
usr/sbin/perdition.imap4
usr/sbin/perdition.pop3
/usr/man/man8/perdition.8
