# debian/perdition.postrm
#! /bin/sh
# postrm script for perdition

#DEBHELPER#

update-rc.d -f perdition remove >/dev/null

exit 0


