# debian/perdition.postrm
#! /bin/sh
# postrm script for perdition

#DEBHELPER#

if [ "$1" = "remove" ]; then
	ldconfig
	update-rc.d -f perdition remove >/dev/null
fi

exit 0


