# debian/perdition.postinst
#! /bin/sh
# postinst script for perdition

update-rc.d perdition defaults >/dev/null

exit 0


