# debian/perdition.postinst
#! /bin/sh
# postinst script for perdition

#DEBHELPER#

if [ "$1" = "configure" ]; then
    ldconfig
    update-rc.d perdition defaults >/dev/null
fi

exit 0


